# systemd/codex-discord-main.service
[Unit]
Description=Codex Discord Main Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/srv/codex-discord/main
EnvironmentFile=/srv/codex-discord/main/.env
ExecStart=/srv/codex-discord/main/.venv/bin/python /srv/codex-discord/main/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
