# systemd/codex-discord-staging.service
[Unit]
Description=Codex Discord Staging Bot
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/qwerasdfzxcl/codex-discord/staging
EnvironmentFile=/home/qwerasdfzxcl/codex-discord/staging/.env
ExecStart=/home/qwerasdfzxcl/codex-discord/staging/.venv/bin/python /home/qwerasdfzxcl/codex-discord/staging/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
