# systemd/codex-discord-staging.service
[Unit]
Description=Codex Discord Staging Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/srv/codex-discord/staging
EnvironmentFile=/srv/codex-discord/staging/.env
ExecStart=/srv/codex-discord/staging/.venv/bin/python /srv/codex-discord/staging/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
