# systemd/codex-discord-main.service
[Unit]
Description=Codex Discord Main Bot
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/codex-discord/main
EnvironmentFile=%h/codex-discord/main/.env
ExecStart=%h/codex-discord/main/.venv/bin/python %h/codex-discord/main/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
