# systemd/codex-discord-staging.service
[Unit]
Description=Codex Discord Staging Bot
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/codex-discord/staging
EnvironmentFile=%h/codex-discord/staging/.env
ExecStart=%h/codex-discord/staging/.venv/bin/python %h/codex-discord/staging/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
