# systemd/codex-discord-main.service
[Unit]
Description=Codex Discord Main Bot
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/qwerasdfzxcl/codex-discord/main
EnvironmentFile=/home/qwerasdfzxcl/codex-discord/main/.env
ExecStart=/home/qwerasdfzxcl/codex-discord/main/.venv/bin/python /home/qwerasdfzxcl/codex-discord/main/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
